# systemd/ingest-tickets.timer
[Unit]
Description=Run ticket ingestion daily at 9:00

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
