# systemd/ingest-tickets.service
[Unit]
Description=MetaShip Ticket Ingestion
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/Support_AI_Research
ExecStart=/usr/bin/node ingest-daily.js
StandardOutput=append:/home/ubuntu/Support_AI_Research/logs/ingest.log
StandardError=append:/home/ubuntu/Support_AI_Research/logs/ingest-error.log

[Install]
WantedBy=multi-user.target
